# moodrain-bookmark-backup/backup.cron
0 * * * * sh /var/script/moodrain-bookmark-backup/cron.sh
